# com.example.daily-research.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yourname.daily-research</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/YOUR_USERNAME/path/to/daily-research/scripts/daily-research.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/path/to/daily-research/logs/launchd-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/path/to/daily-research/logs/launchd-stderr.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/path/to/daily-research</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/YOUR_USERNAME</string>
    <!-- MEM0_API_KEY: Set via 'launchctl setenv MEM0_API_KEY your-key' or in the actual plist.
         Do NOT commit real keys to Git. This file (.example) is a template only. -->
    <key>MEM0_API_KEY</key>
    <string>your-mem0-api-key-here</string>
  </dict>
</dict>
</plist>
